# openclaw-24-04/files/etc/update-motd.d/99-one-click
#!/bin/sh
#
# Configured as part of the DigitalOcean 1-Click Image build process

pub=$(curl -fsS --retry 3 --retry-connrefused --max-time 3 \
  http://169.254.169.254/metadata/v1/interfaces/public/0/ipv4/address 2>/dev/null || true)
myip="${pub:-$(hostname -I | awk '{print$1}')}"
gateway_token=$(grep "^OPENCLAW_GATEWAY_TOKEN=" /opt/openclaw.env 2>/dev/null | cut -d'=' -f2-)

cat <<EOF
********************************************************************************

Welcome to OpenClaw - Your Personal AI Assistant

OpenClaw is a personal AI assistant you run on your own devices. It answers you
on the channels you already use (WhatsApp, Telegram, Slack, Discord, and more).

🌐 Control UI & Gateway Access:
  Dashboard URL: https://$myip
  (Use this URL in the browser — it must match gateway.controlUi.allowedOrigins.)
  (Paste token in dashboard settings if prompted)
  Gateway Token: $gateway_token

  Stuck on "pairing required"? After you connect in the browser, run as root:
    sudo /opt/openclaw-approve-ui-pairing.sh

📝 Configuration:
  Edit settings: /opt/openclaw.env
  User config: /home/openclaw/.openclaw/openclaw.json
  
🔧 Management Commands:
  Restart service: systemctl restart openclaw
  Check status: systemctl status openclaw
  View logs: journalctl -u openclaw -f
  
  Or use helper scripts:
  - /opt/restart-openclaw.sh   (restart with status check)
  - /opt/status-openclaw.sh    (show status and token)
  - /opt/update-openclaw.sh    (update to latest version)
  - /opt/openclaw-cli.sh       (run CLI commands)
  - /opt/openclaw-approve-ui-pairing.sh  (approve Control UI after "pairing required")

🔒 Enable HTTPS (TLS):
  Point your domain to this droplet, then run:
  sudo /opt/setup-openclaw-domain.sh
  
  This configures Caddy with Let's Encrypt SSL certificates.

📱 Configure Messaging Channels:
  1. Edit /opt/openclaw.env with your channel tokens
  2. Or use CLI: /opt/openclaw-cli.sh channels add
  3. Restart: systemctl restart openclaw

📚 Documentation: https://docs.clawd.bot/
🔗 GitHub: https://github.com/openclaw/openclaw
📘 OpenAI vs Codex (API key vs subscription): https://github.com/openclaw/openclaw/blob/main/docs/providers/openai.md

🔧 You can launch OpenClaw TUI using:
  $ /opt/openclaw-tui.sh

********************************************************************************
To delete this message of the day: rm -rf $(readlink -f ${0})
EOF
