# rails-24-04/files/etc/systemd/system/rails-app.service
[Unit]
Description=Rails 8.0.3 Application
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/rails-app
ExecStart=/opt/rails-app/start.sh
ExecStop=/opt/rails-app/stop.sh
TimeoutStartSec=300

[Install]
# Don't auto-start - first boot script will handle initial startup
# WantedBy=multi-user.target