# campfire-24-04/files/etc/update-motd.d/99-one-click
#!/bin/sh
#
# Configured as part of the DigitalOcean 1-Click Image build process

myip=$(hostname -I | awk '{print$1}')
cat <<EOF
********************************************************************************

Campfire, a simple group chat application designed for small teams and organizations. 
Campfire provides real-time messaging with a clean, distraction-free interface 
perfect for focused team communication.

To modify the configuration, access the following file:

  /opt/campfire.env

After making changes, restart Campfire with:

  /opt/restart-campfire.sh

Note: Use the restart script instead of 'docker restart' to ensure 
environment variable changes take effect.

********************************************************************************
To delete this message of the day: rm -rf $(readlink -f ${0})
EOF
