# docker-agent-24-04/files/etc/update-motd.d/99-one-click
#!/bin/sh
#
# Configured as part of the DigitalOcean 1-Click Image build process

myip=$(hostname -I | awk '{print$1}')
cat <<'EOF'
********************************************************************************

Welcome to Docker Agent by Docker Engineering!

Docker Agent is an AI Agent Builder and Runtime. Use the docker-agent CLI to run
and manage agents with tools and YAML-based configuration.

--------------------------------------------------------------------------------
QUICK START
--------------------------------------------------------------------------------

(1) API KEYS

    Export only what your agent YAML needs.

      OPENAI_API_KEY          OpenAI
      ANTHROPIC_API_KEY       Anthropic
      GOOGLE_API_KEY          Google Gemini
      DO_GRADIENT_API_KEY     DigitalOcean Inference (Gradient)

    Example:

      export OPENAI_API_KEY=your_key_here
      export ANTHROPIC_API_KEY=your_key_here
      export GOOGLE_API_KEY=your_key_here
      export DO_GRADIENT_API_KEY=your_key_here

    Gradient keys: https://cloud.digitalocean.com/gen-ai
    All key URLs:  /opt/docker-agent/README.txt

(2) RUN EXAMPLES

    Bundled examples:

      docker-agent run /opt/docker-agent/examples/basic_agent.yaml
      docker-agent run /opt/docker-agent/examples/gradient_agent.yaml
      docker-agent run /opt/docker-agent/examples/dmr.yaml
      docker-agent run /opt/docker-agent/examples/pirate.yaml

    Keys: basic + gradient need OPENAI / DO_GRADIENT; dmr is local; pirate: see YAML.

    Anthropic / Google: set the matching key, then run YAML that uses the same
    token_key. Create a starter with:
      docker-agent new --model anthropic/...
      docker-agent new --model google/...

    One line per provider (same order as keys above):

      OPENAI_API_KEY=k docker-agent run /opt/docker-agent/examples/basic_agent.yaml
      ANTHROPIC_API_KEY=k docker-agent run ./your-agent.yaml
      GOOGLE_API_KEY=k docker-agent run ./your-agent.yaml
      DO_GRADIENT_API_KEY=k docker-agent run /opt/docker-agent/examples/gradient_agent.yaml

(3) CREATE AGENTS

      docker-agent new
      docker-agent new --model openai/gpt-4o-mini

(4) DOCKER HUB

      docker-agent run creek/pirate

--------------------------------------------------------------------------------
RESOURCES
--------------------------------------------------------------------------------

  README & key help:   /opt/docker-agent/README.txt
  Example YAML:        /opt/docker-agent/examples/
  Upstream repo:       https://github.com/docker/docker-agent
  Usage docs:          https://github.com/docker/docker-agent/blob/main/docs/USAGE.md

  docker-agent --help   (CLI reference)

  Variables: OPENAI_API_KEY, ANTHROPIC_API_KEY, GOOGLE_API_KEY, DO_GRADIENT_API_KEY
  Gradient YAML:       base_url https://inference.do-ai.run/v1
                       token_key DO_GRADIENT_API_KEY

********************************************************************************
To delete this message of the day: rm -rf /etc/update-motd.d/99-one-click
EOF
