# docker-agent-24-04/files/etc/update-motd.d/99-one-click
#!/bin/sh
#
# Configured as part of the DigitalOcean 1-Click Image build process

myip=$(hostname -I | awk '{print$1}')
cat <<EOF
********************************************************************************

Welcome to Docker Agent by Docker Engineering!

Docker Agent is an AI Agent Builder and Runtime. Use the docker-agent CLI to run
and manage intelligent AI agents with specialized capabilities and tools.

QUICK START
===========

1. Set your API keys (only what your agent YAML needs):

   export DO_GRADIENT_API_KEY=your_key_here # DigitalOcean Inference (Gradient)
                                            # Create key: https://cloud.digitalocean.com/gen-ai
   export OPENAI_API_KEY=your_key_here      # OpenAI
   export ANTHROPIC_API_KEY=your_key_here   # Anthropic
   export GOOGLE_API_KEY=your_key_here      # Google Gemini

2. Run example agents:

   docker-agent run /opt/docker-agent/examples/gradient_agent.yaml  # needs DO_GRADIENT_API_KEY
   docker-agent run /opt/docker-agent/examples/basic_agent.yaml     # needs OPENAI_API_KEY
   docker-agent run /opt/docker-agent/examples/dmr.yaml             # Docker Model Runner (local)
   docker-agent run /opt/docker-agent/examples/pirate.yaml          # see YAML for provider

   One line (Gradient): DO_GRADIENT_API_KEY=key docker-agent run /opt/docker-agent/examples/gradient_agent.yaml

3. Create new agents:

   docker-agent new                                            # Interactive mode
   docker-agent new --model openai/gpt-4o-mini                # With specific model

4. Pull and run agents from Docker Hub:

   docker-agent run creek/pirate

USEFUL RESOURCES
================

Documentation:       /opt/docker-agent/README.txt
Example agents:      /opt/docker-agent/examples/
GitHub repository:  https://github.com/docker/docker-agent
Usage guide:         https://github.com/docker/docker-agent/blob/main/docs/USAGE.md

Run 'docker-agent --help' for all docker-agent CLI commands and flags.

DigitalOcean Inference (Gradient): export DO_GRADIENT_API_KEY, then
  docker-agent run /opt/docker-agent/examples/gradient_agent.yaml
Custom YAML: base_url https://inference.do-ai.run/v1, token_key DO_GRADIENT_API_KEY

********************************************************************************
To delete this message of the day: rm -rf /etc/update-motd.d/99-one-click
EOF
