# docker-agent-24-04/files/etc/update-motd.d/99-one-click
#!/bin/sh
#
# Configured as part of the DigitalOcean 1-Click Image build process

myip=$(hostname -I | awk '{print$1}')
cat <<EOF
********************************************************************************

Welcome to Docker Agent by Docker Engineering!

Docker Agent is an AI Agent Builder and Runtime. Use the docker-agent CLI to run
and manage intelligent AI agents with specialized capabilities and tools.

QUICK START
===========

1. Set your API keys (choose based on your needs):

   export OPENAI_API_KEY=your_key_here      # For OpenAI models
   export ANTHROPIC_API_KEY=your_key_here   # For Anthropic models
   export GOOGLE_API_KEY=your_key_here      # For Google Gemini models
   export DO_GRADIENT_API_KEY=your_key_here # DigitalOcean Inference (Gradient)

2. Run example agents:

   docker-agent run /opt/docker-agent/examples/basic_agent.yaml     # Requires API key
   docker-agent run /opt/docker-agent/examples/dmr.yaml             # Uses Docker Model Runner
   docker-agent run /opt/docker-agent/examples/pirate.yaml         # Fun pirate assistant

3. Create new agents:

   docker-agent new                                            # Interactive mode
   docker-agent new --model openai/gpt-4o-mini                # With specific model

4. Pull and run agents from Docker Hub:

   docker-agent run creek/pirate

USEFUL RESOURCES
================

Documentation:       /opt/docker-agent/README.txt
Example agents:      /opt/docker-agent/examples/
GitHub repository:  https://github.com/docker/docker-agent
Usage guide:         https://github.com/docker/docker-agent/blob/main/docs/USAGE.md

Run 'docker-agent --help' for all available commands.

DigitalOcean Inference (Gradient): create a key at https://cloud.digitalocean.com/gen-ai
Then use in your agent YAML: provider base_url https://inference.do-ai.run/v1, token_key DO_GRADIENT_API_KEY

********************************************************************************
To delete this message of the day: rm -rf /etc/update-motd.d/99-one-click
EOF
