# hermes-24-04/files/etc/update-motd.d/99-one-click
#!/bin/sh
#
# Configured as part of the DigitalOcean 1-Click Image build process

cat <<'EOF'
********************************************************************************

Welcome to Hermes Agent - The self-improving AI agent from Nous Research

Hermes is installed for the dedicated 'hermes' user. On first SSH login, the
setup wizard runs automatically so you can configure your model provider,
terminal backend, tools, and optional messaging gateway.

Quick start:
  hermes                    Start an interactive Hermes chat
  hermes setup              Re-run setup / reconfigure
  hermes model              Change provider or model
  hermes gateway setup      Configure messaging platforms
  hermes gateway start      Start the messaging gateway
  hermes doctor             Diagnose configuration issues

Files and directories:
  Hermes home:  /home/hermes/.hermes
  Config:       /home/hermes/.hermes/config.yaml
  Env keys:     /home/hermes/.hermes/.env
  Workspace:    /home/hermes/workspace

Helper scripts:
  /opt/hermes/status-hermes.sh
  /opt/hermes/update-hermes.sh
  /opt/hermes/doctor-hermes.sh
  /opt/hermes/hermes-cli.sh

Documentation: https://hermes-agent.nousresearch.com/docs/
GitHub:        https://github.com/NousResearch/hermes-agent

********************************************************************************
To delete this message of the day: rm -rf $(readlink -f ${0})
EOF
