# wordpress-24-04/files/etc/update-motd.d/99-one-click
#!/bin/sh
#
# Configured as part of the DigitalOcean 1-Click Image build process

myip=$(hostname -I | awk '{print$1}')
cat <<EOF
********************************************************************************

Welcome to DigitalOcean's One-Click WordPress Droplet.
To keep this Droplet secure, the UFW firewall is enabled.
All ports are BLOCKED except 22 (SSH), 80 (HTTP), and 443 (HTTPS).

🚀 AUTOMATIC SETUP WITH SSL:
 * The automatic setup wizard will launch now (first login only)
 * Your WordPress site will be secured with HTTPS using Caddy
 * Uses Let's Encrypt short-lived certificates (supports IP addresses!)
 * Setup takes ~2-3 minutes to complete

In a web browser, you can view:
 * Setup pending page (before setup): http://$myip
 * Your WordPress site (after setup): https://$myip
 * WordPress admin panel: https://$myip/wp-admin

On the server:
 * WordPress root: /var/www/html
 * Caddy config: /etc/caddy/Caddyfile
 * SSL certificates: /var/lib/caddy/.local/share/caddy/certificates/
 * Database credentials: /root/.digitalocean_password
   (If using DBaaS: /root/.digitalocean_dbaas_credentials)
 * WP-Fail2Ban plugin: /var/www/html/wp-content/mu-plugins/fail2ban.php

🔧 USEFUL COMMANDS:
 * Add a custom domain: /root/wp_setup_domain.sh
 * WordPress CLI: wp --allow-root --path="/var/www/html" [command]
 * Reload Caddy config: systemctl reload caddy
 * View Caddy logs: journalctl -u caddy -f
 * View PHP logs: tail -f /var/log/php8.3-fpm.log

📚 DOCUMENTATION:
 * WordPress 1-Click Guide: https://do.co/34TfYn8
 * DigitalOcean Community: https://www.digitalocean.com/community/

IMPORTANT:
 * SSL certificates auto-renew every ~6 days (short-lived certificates)
 * Caddy handles renewals automatically in the background
 * For security, xmlrpc endpoint is blocked by default
 * You can switch to a custom domain anytime using /root/wp_setup_domain.sh

********************************************************************************
To delete this message of the day: rm -rf $(readlink -f ${0})
EOF
